# crontab
0 6 * * * cd /app && /usr/local/bin/python3 -u -m generate_epg --source-url "https://epgshare01.online/epgshare01/epg_ripper_AU1.xml.gz" "https://epgshare01.online/epgshare01/epg_ripper_BR2.xml.gz" "https://epgshare01.online/epgshare01/epg_ripper_CA2.xml.gz" "https://epgshare01.online/epgshare01/epg_ripper_CL1.xml.gz" "https://epgshare01.online/epgshare01/epg_ripper_ES1.xml.gz" "https://epgshare01.online/epgshare01/epg_ripper_FR1.xml.gz" "https://epgshare01.online/epgshare01/epg_ripper_IT1.xml.gz" "https://epgshare01.online/epgshare01/epg_ripper_IN1.xml.gz" "https://epgshare01.online/epgshare01/epg_ripper_MX1.xml.gz" "https://epgshare01.online/epgshare01/epg_ripper_NZ1.xml.gz" "https://epgshare01.online/epgshare01/epg_ripper_UK1.xml.gz" "https://epgshare01.online/epgshare01/epg_ripper_US_LOCALS1.xml.gz" "https://epgshare01.online/epgshare01/epg_ripper_US2.xml.gz" >> /proc/1/fd/1 2>&1
